# .claude/skills/WazuhDashboard/wazuh-dashboard.service
[Unit]
Description=Wazuh Security Dashboard
Documentation=https://github.com/cgarrison/PAI
After=network.target

[Service]
Type=simple
User=cgarrison
Group=cgarrison
WorkingDirectory=/home/cgarrison/PAI/.claude/skills/WazuhDashboard
ExecStart=/home/cgarrison/PAI/.claude/skills/WazuhDashboard/manage.sh start
ExecStop=/home/cgarrison/PAI/.claude/skills/WazuhDashboard/manage.sh stop
Restart=on-failure
RestartSec=10

# Environment
Environment="HOME=/home/cgarrison"
Environment="PATH=/home/cgarrison/.bun/bin:/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wazuh-dashboard

[Install]
WantedBy=multi-user.target
