# .claude/skills/WazuhDashboard/wazuh-dashboard.service
[Unit]
Description=Wazuh Security Dashboard
Documentation=https://github.com/YOUR_USERNAME/your-repo
After=network.target

[Service]
Type=simple
User=youruser
Group=youruser
WorkingDirectory=/home/youruser/PAI/.claude/skills/WazuhDashboard
ExecStart=/home/youruser/PAI/.claude/skills/WazuhDashboard/manage.sh start
ExecStop=/home/youruser/PAI/.claude/skills/WazuhDashboard/manage.sh stop
Restart=on-failure
RestartSec=10

# Environment
Environment="HOME=/home/youruser"
Environment="PATH=/home/youruser/.bun/bin:/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wazuh-dashboard

[Install]
WantedBy=multi-user.target
