# .claude/skills/EmailManager/linux-service/emailmanager.service
[Unit]
Description=PAI Email Monitor Service
Documentation=https://github.com/YOUR_USERNAME/your-repo
After=network.target

[Service]
Type=simple
User=youruser
Group=youruser
WorkingDirectory=/home/youruser/.claude/skills/EmailManager
ExecStart=/home/youruser/.bun/bin/bun /home/youruser/.claude/skills/EmailManager/tools/EmailMonitor.ts --daemon
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/youruser/.claude/skills/EmailManager/data/emailmonitor.log
StandardError=append:/home/youruser/.claude/skills/EmailManager/data/emailmonitor.log

# Environment
Environment=HOME=/home/youruser
Environment=PATH=/home/youruser/.bun/bin:/usr/local/bin:/usr/bin:/bin

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/youruser/.claude/skills/EmailManager/data
PrivateTmp=true

[Install]
WantedBy=multi-user.target
